# meta-wifx-openrc/recipes-core/busybox/busybox/busybox-syslogd.initd
#!/sbin/openrc-run

depend() {
	use clock
	need hostname localmount
	after bootmisc
	provide logger
}

command_args="${SYSLOGD_OPTS}"
description="busybox system logger"
command="/sbin/syslogd"

# vim: noet ft=gentoo-init-d
