# meta-wifx-openrc/recipes-core/busybox/busybox/busybox-klogd.initd
#!/sbin/openrc-run

depend() {
	use clock
	need hostname localmount
	after bootmisc
	provide klog
}

command_args="${KLOGD_OPTS}"
description="busybox kernel logger"
command="/sbin/klogd"

# vim: noet ft=gentoo-init-d
