# meta/recipes-extended/logrotate/logrotate/logrotate.cron

#   m    h  dom mon dow user	command
   */10  *     * * *    root   /etc/cron.script/logrotate.sh
