# layers/meta-openrc-app/recipes-connectivity/openssh/openssh/ssh-key.initd
#!/sbin/openrc-run

description="SSH daemon key generation"

depend() {
    keyword -shutdown -stop
}

start() {
    ssh-keygen -A
    eend $?
}