# meta/recipes-core/base-files/base-files/update-motd.d/00-header
#!/bin/sh

cat /etc/update-motd.d/static/00-header
echo ""