# meta/recipes-core/base-files/base-files/update-motd.d/90-release-notes
#!/bin/sh

cat /etc/update-motd.d/static/90-release-notes
