# meta/recipes-core/base-files/base-files/update-motd.d/99-reset-style
#!/bin/sh

# Ensure background and foreground styles are set by default again
echo "[0;39m[0;49m"
