# layers/target/meta-wifx-machine/recipes-core/led-service/led-service/led-service.initd
#!/sbin/openrc-run

export EINFO_QUIET="${quiet:-yes}"

depend() {
        after local
        keyword -prefix
}

start() {
        if [ -e /sys/class/leds/service/enable ]; then
                echo "1" > /sys/class/leds/service/enable
        fi
}